# ums.timer
[Unit]
Description=Start UMS (Beta) delayed on boot

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
